# failover_oracle/ansible/roles/failover_cf/templates/failback.service.j2
[Unit]
Description=Cloudflare DNS daily failback to PRIMARY (probe via tunnel)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment="TZ={{ _time_zone }}"
ExecStart={{ _failover_dir }}/failback.sh
